# audit-monitor.service
[Unit]
Description=Monitor sensitive file access using auditd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/audit-monitor.sh
